# infra/systemd/lcdmonitor.system.service
[Unit]
Description=LCD Monitor Server (system mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as a non-root user; adjust to your deployment user
User=lcdmon
Group=dialout
WorkingDirectory=/opt/lcdmonitor/server
EnvironmentFile=/etc/default/lcdmonitor
# Use a shell wrapper so the EnvironmentFile variables expand correctly
ExecStart=/bin/sh -c 'exec "$${LCDMONITOR_VENV}/bin/python" -m src.main --config "$${LCDMONITOR_CONFIG}" --exec-driver shell'
Restart=on-failure
RestartSec=5

# Hardening (tune as needed)
# NoNewPrivileges must remain disabled so sudo can elevate for shell execs
NoNewPrivileges=no
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectClock=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
SystemCallArchitectures=native
MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target
